# etc/systemd/system/bw-auto-confirm.service
[Unit]
Description=Bitwarden auto-confirm organization members

[Service]
Type=oneshot
ExecStart=/bin/bash /root/bwAutoConfirm.sh
